# target/PolarfireBridge/ldmx-dma-read.service
[Unit]
Description=ldmx-dma-read

[Service]
ExecStart=/usr/local/sbin/dma-read

[Install]
WantedBy=multi-user.target
